# restic-backup.timer


[Unit]
Description=Timer for Running Restic Backup


[Timer]
OnBootSec=5min
Unit=restic-backup.service


[Install]
WantedBy=timers.target
